# etc/init/turnserver.conf
# matrix-synapse service

description     "matrix-synapse"
author          "Martin Giess <martin.giess@emnify.com>"

start on runlevel [2345]

stop on starting rc RUNLEVEL=[016]

respawn
respawn limit 2 5

setuid {{ username }}

exec turnserver
